# phonebook/files/fetch_phonebook.cron
#!/bin/sh

URL="http://hb9edi-apu-1.local.mesh:8080/filerepo/Phonebook/AREDN_Phonebook.csv"
FILE="/www/phonebook.csv"
if [ -f "$FILE" ]; then
  if [[ $(find "$FILE" -mmin +1400 -print) ]]; then
    curl -o $FILE $URL
  fi
else
  curl -o $FILE $URL
  /etc/cron.hourly/fetch_phonebook
fi