# phonebook/files/update_phonebook.cron
#!/bin/sh
/usr/bin/phonebook --conf=/etc/config --server=false
