# phonebook/files/update_phonebook.cron
#!/bin/sh
/usr/bin/phonebook --conf=/etc/phonebook.conf --server=false
