# phonebook-deb/DEBIAN/prerm
#!/bin/sh
set -e

#DEBHELPER#

echo "Stop phonebook service"
systemctl stop phonebook
systemctl disable phonebook

exit 0
